# add/etc/service/ngrok/run
#!/bin/sh

/sbin/setuser ngrok ngrok start --all
